# landing/nexa-trials/marketing/lab-vps/nexa-labs.service
[Unit]
Description=Nexa lab demos (Docker Compose)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/nexa-labs/landing/nexa-trials/marketing/lab-vps
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
